# src/bash/led-d.service
[Unit]
Description=Led matrix server

[Service]
User=pi
ExecStart=/usr/bin/led-d

[Install]
WantedBy=multi-user.target
